# deamon/daily.cron
0 7 * * * python3 /app/dailyd.py >>/app/out/dailyd.out 2>&1
